# local-service-deploy/systemd/system/docker-cleanup.timer
[Unit]
Description=Docker cleanup timer

[Timer]
OnUnitInactiveSec=48h

[Install]
WantedBy=timers.target
